# templates/default/npmlazy.conf.erb
# Ubuntu upstart file at /etc/init/npmlazy.conf
# Automatically Generated by Chef

description "NPM lazy cache Node.JS Server"

start on (starting network-interface
          or starting network-manager
          or starting networking)

stop on runlevel [!023456]

setuid <%= node['npmlazy']['user'] %>
setgid <%= node['npmlazy']['user'] %>

respawn

script
  exec npm_lazy -c <%= File.join(node['npmlazy']['conf']['confdir'], 'npm_lazy_config.js') %>
end script
