# templates/default/sinopia.service.erb
# Automatically Generated by Chef

[Unit]
Description=Sinopia NPM cache & private repo server

[Service]
ExecStart=/usr/bin/sinopia -c <%= File.join(node['sinopia']['confdir'], 'config.yaml') %>
Type=simple
User=sinopia
Group=sinopia

[Install]
WantedBy=multi-user.target
